# fabfile/templates/shipyard-agent.conf
[program:shipyard-agent]
directory=/tmp
user=root
command=/usr/local/bin/shipyard-agent
    -url %(url)s
    -key %(key)s
autostart=true
autorestart=true
